# eagler-files/rootfs/etc/update-motd.d/10-help-text
version https://git-lfs.github.com/spec/v1
oid sha256:dc3a00a8bbe51dad1e958d6c5f73a59a6ed61a9e37acec36a78a1e817e5cadd4
size 1151
