# eagler-files/rootfs/etc/update-motd.d/50-motd-news
version https://git-lfs.github.com/spec/v1
oid sha256:12d09f922c013d07b62cca4fda15a567d94133b196c23193784f1725562ad754
size 5023
