# eagler-files/rootfs/etc/update-motd.d/60-unminimize
version https://git-lfs.github.com/spec/v1
oid sha256:e1dc5ce642b5de0eb7c0c20d3362fd7c9d1748befed8695586f59626108abbd9
size 356
